# data/fwupd-offline-update.service.in
[Unit]
Description=Updates device firmware whilst offline
Documentation=man:fwupdmgr
ConditionPathExists=@localstatedir@/lib/fwupd/pending.db
DefaultDependencies=false
Requires=fwupd.sevice

[Service]
Type=oneshot
ExecStart=@libexecdir@/fwupd/fwupdoffline

[Install]
WantedBy=shutdown.target reboot.target halt.target
